# deploy/moestuin-webcam.service
[Unit]
Description=Moestuin webcam capture (HLS)
After=network-online.target

[Service]
Type=simple
User=moestuin
Group=moestuin
SupplementaryGroups=video
WorkingDirectory=/var/lib/moestuin/webcam
ExecStart=/usr/local/bin/moestuin-webcam.sh
Restart=always
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
