# deploy/moestuin-backup.service
[Unit]
Description=Moestuin nightly backup + webcam rollup

[Service]
Type=oneshot
User=moestuin
Group=moestuin
ExecStart=/usr/local/bin/moestuin-backup.sh
